# xleax_profile/airootfs/etc/systemd/system/btrfs-scrub@.timer.d/weekly.conf
[Timer]
OnCalendar=
OnCalendar=weekly
Persistent=true