# xleax_profile/airootfs/etc/xdg/autostart/xleax-proton-setup.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=xleaxOS Proton-GE Setup
Comment=Install latest Proton-GE for Steam on first session
Exec=/usr/local/bin/xleax-proton-setup
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true