# xleax_profile/airootfs/etc/xdg/autostart/xleax-hud.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=EpsteinOS X-HUD
Comment=Launch the EpsteinOS performance HUD overlay
Exec=/usr/local/bin/xleax-hud-start
Terminal=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
