# xleax_profile/airootfs/etc/udev/rules.d/70-mouse-polling.rules
ACTION=="add|change", SUBSYSTEM=="input", KERNEL=="event*", ENV{ID_INPUT_MOUSE}=="1", RUN+="/usr/bin/xleax-mouse-polling %p"
