# xleax_profile/airootfs/etc/systemd/system/xleax-gpu-boost.service
[Unit]
Description=EpsteinOS GPU Performance Boost
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=oneshot
ExecStart=/usr/bin/xleax-gpu-boost
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
